# packages/Debian/common/DEBIAN/prerm
#!/usr/bin/env bash

##--------------------------------------------------------------------
## Copyright (c) 2018 OSIsoft, LLC
##
## Licensed under the Apache License, Version 2.0 (the "License");
## you may not use this file except in compliance with the License.
## You may obtain a copy of the License at
##
##     http://www.apache.org/licenses/LICENSE-2.0
##
## Unless required by applicable law or agreed to in writing, software
## distributed under the License is distributed on an "AS IS" BASIS,
## WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
## See the License for the specific language governing permissions and
## limitations under the License.
##--------------------------------------------------------------------

##--------------------------------------------------------------------
##
## @prerm DEBIAN/prerm
## This script is used to execute before the removal of files associated with the package.
##
## Author: Ivan Zoratti, Ashwin Gopalakrishnan
##
##----------------------------------------------------------------------------------------

set -e

PKG_NAME="fledge"

get_fledge_script () {
    fledge_script=$(dpkg -L $PKG_NAME | grep 'fledge/bin/fledge$')
    echo $fledge_script
}

stop_fledge_service () {
    systemctl stop fledge
}

is_fledge_running () {
    set +e
    fledge_script=$(get_fledge_script)
    fledge_status_output=$($fledge_script status 2>&1 | grep 'Fledge Uptime')
    rc=$((!$?))
    echo $rc
    set -e
}

kill_fledge () {
    set +e
    fledge_script=$(get_fledge_script)
    fledge_status_output=$($fledge_script kill 2>&1)
    set -e
}

disable_fledge_service () {
    systemctl disable fledge
}

remove_fledge_service_file () {
    rm -rf /etc/init.d/fledge
}

reset_systemctl () {
    systemctl daemon-reload
    systemctl reset-failed
}

# main

IS_FLEDGE_RUNNING=$(is_fledge_running)
if [ "$IS_FLEDGE_RUNNING" -eq "1" ]
then
    echo "Fledge is currently running."
    echo "Stop Fledge service."
    stop_fledge_service
    echo "Kill Fledge."
    kill_fledge
fi

echo "Disable Fledge service."
disable_fledge_service
echo "Remove Fledge service script"
remove_fledge_service_file
echo "Reset systemctl"
reset_systemctl
